# build-files/etc/supervisord.d/crond.ini
[program:crond]
command=/bin/bash -c "source /etc/sysconfig/crond && exec /usr/sbin/crond -n $CRONDARGS"
process_name=%(program_name)s
numprocs=1
directory=/tmp
umask=022
autostart=%(ENV_SUPERVISOR_AUTOSTART_CROND)s
autorestart=%(ENV_SUPERVISOR_AUTOSTART_CROND)s
startretries=2
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=root