# build-files/etc/supervisord.d/nginx.ini
[program:nginx]
command=nginx -c /etc/nginx/nginx.conf  -g 'daemon off;'
process_name=%(program_name)s
numprocs=1
autostart=%(ENV_SUPERVISOR_AUTOSTART_NGINX)s
autorestart=%(ENV_SUPERVISOR_AUTORESTART_NGINX)s
startretries=2
stdout_logfile=%(ENV_NGINX_STDOUT_LOGFILE)s
stdout_logfile_maxbytes=0
stderr_logfile=%(ENV_NGINX_STDERR_LOGFILE)s
stderr_logfile_maxbytes=0