# build-files/etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=php-fpm -R -F -c /etc/php-fpm.conf
process_name=%(program_name)s
numprocs=1
autostart=%(ENV_SUPERVISOR_AUTOSTART_PHPFPM)s
autorestart=%(ENV_SUPERVISOR_AUTORESTART_PHPFPM)s
startretries=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=998